# .config/autostart/picom.desktop
[Desktop Entry]
Exec=picom -b
Name=picom
Type=Application
Version=1.0
