# .config/nvim/template/.service
[Service]
ExecStart=/usr/bin/%FILE
ExecReload=/bin/kill -SIGHUP $MAINPID
