# .config/nvim/templates/%.service
[Service]
ExecStart=/usr/bin/{{ expand('%:t:r') }}
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy=default.target
